# overlay/etc/systemd/system/chrony.service.d/pps0-dependency.conf
[Unit]
After=dev-pps0.device
Wants=dev-pps0.device
